# services/bun-dns.timer
[Unit]
Description=Timer for Porkbun DNS Updater

[Timer]
Unit=bun-dns.service
OnBootSec=30s
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
